# wifihotspot.service
[Unit]
Description=wifihotspot

[Service]
Type=simple
ExecStart=/bin/bash /home/openhabian/wifi.sh

[Install]
WantedBy=multi-user.target
